# deploy/fund-watch.timer
[Unit]
Description=基金监控 — 每日晚报（工作日 15:30）

[Timer]
OnCalendar=Mon..Fri 15:30:00
Persistent=false
RandomizedDelaySec=10

[Install]
WantedBy=timers.target
